# data/export/upstartrhel6/process.conf.erb
start on runlevel [2345]
respawn

<% engine.env.each do |name,value| -%>
env <%= name.upcase %>='<%= value.gsub(/'/, "'\"'\"'") %>'
<% end -%>

exec su - <%= user %> -c 'cd <%= engine.root %>; ; export PORT=<%= port %>; <%= process.command %>' 
